# hermes/scripts/qnoe-hermes.service.d/50-b7-readonly.conf
# B7 read-only enforcement (2026-07-14) — see B7_SANDBOX_HANDOFF.md, redteam R4.
# Namespace-level mounts: binds write_file/patch AND every terminal child process.
# ROLLBACK: overwrite this file with an empty one (sudo cp /dev/null <this file>),
#           then: sudo systemctl daemon-reload && sudo systemctl restart qnoe-hermes
[Service]
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
# Everything under /opt/qnoe-agent, the lab share, and yzamir's home (qcodes
# registry) is read-only for the gateway and all its children...
ReadOnlyPaths=/opt/qnoe-agent /ICFO /mnt/noe /home/yzamir
# ...except the state the gateway legitimately writes: Mem0/SQLite, logs,
# Hermes session/profile state.
ReadWritePaths=/opt/qnoe-agent/memory /opt/qnoe-agent/logs /opt/qnoe-agent/hermes
# Secrets are invisible inside the namespace; teams.env is delivered by systemd
# (read as root, outside the namespace) via $CREDENTIALS_DIRECTORY instead.
InaccessiblePaths=/opt/qnoe-agent/secrets
LoadCredential=teams.env:/opt/qnoe-agent/secrets/teams.env
